# base-station-config/50-logitech-f710.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c21f", RUN+="/bin/bash /etc/udev/rules.d/f710-setup.sh"
