# linux-config/90-gamepad.rules
ATTRS{idVendor}=="20d6", ATTRS{idProduct}=="2019", RUN+="/usr/local/bin/map_gamepad.sh"
